# navigator/home/navigator/.config/autostart/no-screen-blank.desktop
[Desktop Entry]
Type=Application
Name=Disable Screen Blanking
Exec=xset s off -dpms
NoDisplay=true
X-XFCE-Autostart-Phase=Applications
